# systemd/user/wallpaper.path
[Unit]
Description=Wallpaper monitor

[Path]
PathExists=/tmp/update_wallpaper

[Install]
WantedBy=default.target